# dnscrypt-proxy
#!/sbin/openrc-run

depend() {
        before netmount nfsmount
        use net
}

start() {
        ebegin "Starting dnscrypt-proxy"
        /sbin/start-stop-daemon -b --exec --start -m --pidfile /run/dnscrypt-proxy.pid /usr/sbin/dnscrypt-proxy -- -config /etc/dnscrypt-proxy/dnscrypt-proxy.toml
        eend $?
}

stop() {
        ebegin "Stopping dnscrypt-proxy"
        kill -TERM -`cat /run/dnscrypt-proxy.pid`
        rm -f /run/dnscrypt-proxy.pid
        eend $?
}
